# rootfs/etc/s6-overlay/s6-rc.d/init-mailcue/up
/bin/sh -e /etc/s6-overlay/scripts/init-mailcue.sh
