# rootfs/etc/s6-overlay/s6-rc.d/acme-setup/up
with-contenv /bin/sh -e /etc/s6-overlay/scripts/acme-setup.sh
