# rootfs/etc/s6-overlay/s6-rc.d/spamass-milter/run
#!/command/execlineb -P
with-contenv
/usr/sbin/spamass-milter -p /var/run/spamass/spamass.sock -f -u debian-spamd -- -s 127.0.0.1:783
