# rootfs/etc/s6-overlay/s6-rc.d/spamass-milter/run
#!/command/execlineb -P
with-contenv
/usr/sbin/spamass-milter -p inet:10028@127.0.0.1 -f -u debian-spamd -m -- -s 127.0.0.1:783
