# rootfs/etc/s6-overlay/s6-rc.d/spamd/run
#!/bin/sh
exec 2>&1

MAX_CHILDREN="${MAILCUE_SPAMD_MAX_CHILDREN:-5}"
exec /usr/sbin/spamd --nouser-config --max-children="${MAX_CHILDREN}" --listen=127.0.0.1:783
