# rootfs/etc/s6-overlay/s6-rc.d/acme-renew/run
#!/command/execlineb -P
with-contenv
/bin/sh /etc/s6-overlay/scripts/acme-renew.sh
