# rootfs/etc/s6-overlay/s6-rc.d/opendkim/run
#!/bin/sh
# =============================================================================
# OpenDKIM — DKIM signing milter (foreground)
# =============================================================================
exec 2>&1

# Ensure the run directory exists
mkdir -p /var/run/opendkim
chown opendkim:opendkim /var/run/opendkim

# Add postfix to opendkim group so it can connect to the milter socket
adduser postfix opendkim 2>/dev/null || true

# -f = foreground mode
exec /usr/sbin/opendkim -f -x /etc/opendkim/opendkim.conf
